# api/cli/devOps/files/services.d/elixir-app/run
#!/usr/bin/with-contenv sh

exec /usr/bin/elixir-app.sh
